# launchd/com.steve.homemonitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.steve.homemonitor</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>REPLACE_WITH_PROJECT_DIR/scripts/nest_monitor.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>REPLACE_WITH_PROJECT_DIR</string>

    <!-- Poll every 15 minutes (900 seconds) -->
    <key>StartInterval</key>
    <integer>900</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>REPLACE_WITH_PROJECT_DIR/logs/monitor.log</string>

    <key>StandardErrorPath</key>
    <string>REPLACE_WITH_PROJECT_DIR/logs/monitor.err</string>
</dict>
</plist>
