# launchd/com.steve.homeserver.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.steve.homeserver</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>-m</string>
        <string>http.server</string>
        <string>8080</string>
    </array>

    <!-- Run in the project directory so relative paths resolve correctly -->
    <key>WorkingDirectory</key>
    <string>REPLACE_WITH_PROJECT_DIR</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>REPLACE_WITH_PROJECT_DIR/logs/server.log</string>

    <key>StandardErrorPath</key>
    <string>REPLACE_WITH_PROJECT_DIR/logs/server.err</string>
</dict>
</plist>
